# os/systemd/proof-of-life-node.service
[Unit]
Description=Proof of Life Chain Node
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=prooflife
Group=prooflife
WorkingDirectory=/opt/proof-of-life
Environment=NODE_ENV=production
Environment=PROOF_OF_LIFE_CONFIG=/etc/proof-of-life/node.config.json
ExecStart=/usr/bin/node /opt/proof-of-life/chain/dist/node/devnet.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/proof-of-life /var/log/proof-of-life
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
